# deploy/deerflow.service.d/override.conf
[Service]
# DeerFlow Skill 路径注入配置
# 文件: deploy/deerflow.service.d/override.conf
# 部署命令: sudo cp deploy/deerflow.service.d/override.conf /etc/systemd/system/deerflow.service.d/
#           sudo systemctl daemon-reload && sudo systemctl restart deerflow

# 加载环境变量文件（推荐方式）
EnvironmentFile=/etc/default/deerflow-env

# 直接注入环境变量（备选方案，当 EnvironmentFile 不可用时）
# Environment="WATER_RESOURCES_ROOT=/mnt/skills"
# Environment="SL323_DB_HOST=192.168.100.103"
# Environment="SL323_DB_PORT=3306"
# Environment="SL323_DB_USER=root"
# Environment="SL323_DB_PASSWORD=your-password"
